# supervisord/conf.d/cobblerd.conf
[program:cobblerd]
command=/usr/bin/cobblerd -F
redirect_stderr=true
