# supervisord/conf.d/forklift.conf
[program:forklift]
command=/opt/cobblerimporter/cobblerforklift.sh
redirect_stderr=true
restart=false
