# supervisord/conf.d/tftpd.conf
[program:tftpd]
command=/usr/sbin/in.tftpd -v -v -v -4 -L -s /var/lib/tftpboot
redirect_stderr=true
